# crontab
*/5 * * * * /bin/echo "Cron is running at $(date)" >> /var/log/cron.log 2>&1
*/5 * * * * /usr/bin/python3 /app/scraper.py >> /var/log/cron.log 2>&1





