# crontab
*/5 * * * * /usr/local/bin/python3 /app/scraper.py >> /var/log/cron.log 2>&1

